# src/init/dmon-es.conf
description "ElasticSearch service"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [016]

respawn
respawn limit 10 30

limit nofile 32000 32000

setuid ubuntu
setgid ubuntu

env heap_size="4g"

script

  ES_HEAP_SIZE="$heap_size"

  export ES_HEAP_SIZE

  exec /opt/elasticsearch/bin/elasticsearch

end script
